# misc/systemd/nix-daemon.service.in
[Unit]
Description=Nix Daemon
Documentation=man:nix-daemon https://nixos.org/manual
RequiresMountsFor=@storedir@
RequiresMountsFor=@localstatedir@
RequiresMountsFor=@localstatedir@/nix/db
ConditionPathIsReadWrite=@localstatedir@/nix/daemon-socket

[Service]
ExecStart=@@bindir@/nix-daemon nix-daemon --daemon
KillMode=process

[Install]
WantedBy=multi-user.target
