# misc/systemd/nix-daemon.service.in
[Unit]
Description=Nix Daemon
RequiresMountsFor=@storedir@
RequiresMountsFor=@localstatedir@
ConditionPathIsReadWrite=@localstatedir@/nix/daemon-socket

[Service]
ExecStart=@@bindir@/nix-daemon nix-daemon --daemon
KillMode=process

[Install]
WantedBy=multi-user.target
