# misc/systemd/nix-daemon.socket.in
[Unit]
Description=Nix Daemon Socket
Before=multi-user.target
RequiresMountsFor=@storedir@
# Skip activating nix-daemon.socket in case the socket already exists,
# which is the case if we invoke a container and bind-mount the daemon
# socket inside.
ConditionPathExists=!@localstatedir@/nix/daemon-socket/socket

[Socket]
ListenStream=@localstatedir@/nix/daemon-socket/socket

[Install]
WantedBy=sockets.target
