# nix/misc/upstart/nix-daemon.conf.in
description "Nix Daemon"
start on filesystem
stop on shutdown
respawn
exec @bindir@/nix-daemon --daemon
