# misc/openrc/nix-daemon.in
#!/sbin/openrc-run
description="Nix Daemon"
command=@bindir@/nix-daemon
command_background="yes"
pidfile="/run/$RC_SVCNAME.pid"
